# meta-facebook/meta-minerva-cmc/recipes-phosphor/gpio/phosphor-gpio-monitor/blade-not-present@.service
[Unit]
Description=record blade:%i not present

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/record_blade_present_log.sh 1 %i