# meta-facebook/meta-minerva/recipes-phosphor/console/obmc-console/obmc-console-ttyS0-ssh-mtia-blade6.socket
[Unit]
Description=Phosphor BMC Console SSH socket for MTIA Blade 6
Wants=obmc-console@ttyS0.service

[Socket]
ListenStream=2206
Accept=yes

[Install]
WantedBy=sockets.target
