# meta-facebook/meta-minerva/recipes-phosphor/console/obmc-console/obmc-console-ttyS0-ssh-mtia-blade-tmplate.socket
[Unit]
Description=Phosphor BMC Console SSH socket for MTIA Blade {blade_num}
Wants=obmc-console@ttyS0.service

[Socket]
ListenStream={port_num}
Accept=yes

[Install]
WantedBy=sockets.target
