# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/blade-prsnt-deassert@.service
[Unit]
Description=Blade %i Present Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh blade_prsnt 1 %i 