# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/blade-prsnt-assert@.service
[Unit]
Description=Blade %i Present Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh blade_prsnt 0 %i 