# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/blade-pwrok-assert@.service
[Unit]
Description=Blade %i PWROK Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh blade_pwrok 0 %i 