# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/blade-pwrok-deassert@.service
[Unit]
Description=Blade %i PWROK Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh blade_pwrok 1 %i 