# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/cable-prsnt-assert@.service
[Unit]
Description=%i Cable Present Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh cable_prsnt 0 %i