# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/cable-prsnt-deassert@.service
[Unit]
Description=%i Cable Present Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh cable_prsnt 1 %i