# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/leak-detect-assert@.service
[Unit]
Description=%i Leak Detect Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh leak_detect 0 %i 