# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/leak-detect-deassert@.service
[Unit]
Description=%i Leak Detect Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh leak_detect 1 %i 