# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/mgmt-sfp-prsnt-assert.service
[Unit]
Description=MGMT-SFP Present Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh mgmt_sfp_prsnt 0 cmm_blade
