# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/mgmt-sfp-present.service
[Unit]
Description= record mgmt SFP is present

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/mgmt_sfp_present.sh 0