# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/mgmt-sfp-prsnt-deassert.service
[Unit]
Description=MGMT-SFP Present Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh mgmt_sfp_prsnt 1 cmm_blade