# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/toggle-switch-assert.service
[Unit]
Description=CR Toggle Switch Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh toggle_switch 0 cmm_blade