# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/toggle-switch-deassert.service
[Unit]
Description=CR Toggle Switch Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging.sh toggle_switch 1 cmm_blade