# meta-facebook/meta-minerva/recipes-phosphor/state/phosphor-state-manager/chassis-powercycle@.service
[Unit]
Description=Chassis%i Power Cycle 

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-powercycle

[Install]
WantedBy=obmc-chassis-powercycle@%i.target